# config/systemd/user/chrome.service
[Unit]
Descripton=Goole Chrome service

[Service]
ExecStart=/usr/bin/google-chrome-stable

[Install]
WantedBy=default.target
